# roles/smartpanda/upstart/smartpanda.conf
start on runlevel [2345]
stop on runlevel [016]

respawn
exec /usr/bin/nodejs /panda/smartpanda-app/smartpanda.js #&> /panda/smartpanda.log